# nahraj-fotky/psychosevdava.service
[Unit]
Description=Gunicorn instance to serve photo upload for http://psychosevdava.cz/nahraj-fotky/
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/psychosevdava/nahraj-fotky/ 
ExecStart=/var/www/psychosevdava/nahraj-fotky/gunicorn_start.sh # IMPORTANT: CHANGE THIS TO YOUR ACTUAL PATH
Restart=always
PrivateTmp=true

# Environment variables (optional, but good for Flask/Gunicorn)
# Environment="FLASK_ENV=production"
# Environment="SECRET_KEY=your_super_secret_key" # Replace with a real key!

# Standard output and error to syslog
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=my_flask_app

[Install]
WantedBy=multi-user.target
